# ubuntu/shadowsocks.conf
[program:shadowsocks]
command=ssserver -c /etc/shadowsocks_1.json
autostart=true
autorestart=true
user=root
